# api/bambuddy-bug-relay.service
# Copy to /etc/systemd/system/bambuddy-bug-relay.service
# Then: systemctl daemon-reload && systemctl enable --now bambuddy-bug-relay
#
# Create /etc/bambuddy-relay/env with:
#   BUG_REPORT_GITHUB_TOKEN=ghp_your_token_here
# Then: chmod 600 /etc/bambuddy-relay/env

[Unit]
Description=Bambuddy Bug Report Relay
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/bambuddy-relay/api
EnvironmentFile=/etc/bambuddy-relay/env
ExecStart=/opt/bambuddy-relay/venv/bin/gunicorn --bind 127.0.0.1:5050 --workers 2 wsgi:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
