# api/bambuddy-bug-relay.service
# Copy to /etc/systemd/system/bambuddy-bug-relay.service
# Then: systemctl daemon-reload && systemctl enable --now bambuddy-bug-relay
#
# Create /etc/bambuddy-relay/env with all required variables (see below).
# Then: chmod 600 /etc/bambuddy-relay/env
#
# Required env vars in /etc/bambuddy-relay/env:
#   BUG_REPORT_GITHUB_TOKEN=ghp_...
#   SMTP_HOST=smtp.example.com
#   SMTP_PORT=587
#   SMTP_USERNAME=user@example.com
#   SMTP_PASSWORD=secret
#   SMTP_FROM_EMAIL=noreply@bambuddy.cool
#   SMTP_USE_TLS=starttls          # "starttls", "ssl", or "none"
#   MAINTAINER_EMAIL=mz@v8w.de

[Unit]
Description=Bambuddy Bug Report Relay
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/bambuddy-relay/api
EnvironmentFile=/etc/bambuddy-relay/env
ExecStart=/opt/bambuddy-relay/venv/bin/gunicorn --bind 127.0.0.1:5050 --workers 2 wsgi:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
