# api/bambuddy-bug-relay.service
# Copy to /etc/systemd/system/bambuddy-bug-relay.service
# Then: systemctl daemon-reload && systemctl enable --now bambuddy-bug-relay
#
# Set BUG_REPORT_GITHUB_TOKEN to a GitHub PAT with repo scope for maziggy/bambuddy

[Unit]
Description=Bambuddy Bug Report Relay
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/bambuddy-relay/api
Environment=BUG_REPORT_GITHUB_TOKEN=<configure-on-server>
ExecStart=/opt/bambuddy-relay/venv/bin/gunicorn --bind 127.0.0.1:5050 --workers 2 wsgi:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
